# docker/service/nginx/run
#!/bin/sh
exec /usr/sbin/nginx -g 'daemon off;'